# target/linux/ar71xx/base-files/etc/uci-defaults/05_ar750ac_mac
#!/bin/sh

. /lib/ar71xx.sh
. /lib/functions.sh
. /lib/functions/system.sh



set_5g_macaddr()
{
	config_load wireless
	config_foreach set_macaddr wifi-iface
}

set_macaddr()
{
	local device path

	config_load wireless

	config_get device "$1" device
	config_get path $device path

	if [ -n "`echo $path | grep "pci0000"`" ] ; then
		uci -q set wireless.$1.macaddr="$(macaddr_add $(cat /sys/class/net/eth0/address) +1)"
		uci -q commit wireless
	fi
}


board=$(ar71xx_board_name)

case "$board" in
		gl-ar750ac)
			set_5g_macaddr
			;;
		*)
			;;
	esac
exit 0
