# package/ar750-sta-watch/files/ar750-sta-watch
#!/bin/sh /etc/rc.common


START=36
start()
{
     sta-watch &
}

stop()
{
    killall sta-watch
}

restart()
{
	stop
	sleep 1
	start
}
